# example/docker-plugin-cinder.socket
[Socket]
ListenStream=/run/docker/plugins/cinder.sock

[Install]
WantedBy=docker.service
